# .config/autostart/welcome.desktop
[Desktop Entry]
# Put in $HOME/.config/autostart to launch on login
Type=Application
Name=Welcome to Pirate Box
Path=$HOME/bin
Exec=welcome.bash
